# docker/worker/supervisord.conf
[program:sino-worker]
command=php artisan kafka:name
;user=root
directory=/var/www/html/
process_name=sino-worker-%(process_num)
numprocs=31
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/sino.out.log
stderr_logfile=/var/log/supervisor/sino.error.log